# deploy/systemd/automation-worker.timer
# deploy/systemd/automation-worker.timer — schedules automation worker batches
#
# Usage:
#   1. Run `bun run init:site` to replace <project> placeholders.
#   2. Copy to ~/.config/systemd/user/<project>-automation-worker.timer
#   3. systemctl --user daemon-reload
#   4. systemctl --user enable --now <project>-automation-worker.timer
#
# Status / next-fire:
#   systemctl --user list-timers | grep <project>-automation-worker

[Unit]
Description=Automation worker timer — <project>

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
RandomizedDelaySec=10
Persistent=true
Unit=<project>-automation-worker.service

[Install]
WantedBy=timers.target
