# deploy/systemd/backup-base.service
# deploy/systemd/backup-base.service — daily WAL-G base backup for <project>.
#
# Pairs with backup-base.timer. Reads R2/WAL-G credentials from the project
# env file via the postgres container's own environment — this unit only
# runs the backup-base.sh wrapper which `podman exec`s into the container.
#
# Install per project:
#   cp deploy/systemd/backup-base.service ~/.config/systemd/user/<project>-backup-base.service
#   cp deploy/systemd/backup-base.timer   ~/.config/systemd/user/<project>-backup-base.timer
#   systemctl --user daemon-reload
#   systemctl --user enable --now <project>-backup-base.timer

[Unit]
Description=WAL-G base backup — <project>
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/<project>
EnvironmentFile=%h/secrets/<project>.prod.env
ExecStart=%h/<project>/scripts/backup-base.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
