# deploy/systemd/restore-drill.timer
# deploy/systemd/restore-drill.timer — weekly non-destructive restore drill.
#
# OnCalendar fires weekly with a random spread so multiple sites on the same
# host do not all restore from R2 simultaneously.
#
# Pairs with restore-drill.service. See deploy/systemd/restore-drill.service.

[Unit]
Description=Weekly non-destructive restore drill timer — <project>

[Timer]
OnCalendar=Sun *-*-* 03:00:00 UTC
Persistent=true
RandomizedDelaySec=900
Unit=<project>-restore-drill.service

[Install]
WantedBy=timers.target
