# deploy/systemd/restore-drill.service
# deploy/systemd/restore-drill.service — weekly non-destructive restore drill for <project>.
#
# Pairs with restore-drill.timer. The drill restores the latest WAL-G base
# backup into a scratch container, replays WAL to the target time, runs a
# read-only sanity query, tears down, and writes evidence to the ops-status
# ledger.
#
# Install per project:
#   cp deploy/systemd/restore-drill.service ~/.config/systemd/user/<project>-restore-drill.service
#   cp deploy/systemd/restore-drill.timer   ~/.config/systemd/user/<project>-restore-drill.timer
#   systemctl --user daemon-reload
#   systemctl --user enable --now <project>-restore-drill.timer

[Unit]
Description=Weekly non-destructive restore drill — <project>
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/<project>
EnvironmentFile=%h/secrets/<project>.prod.env
ExecStart=%h/.bun/bin/bun run backup:restore:drill
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
